# kylin/supervisor.d/kylin.conf
[program:kylin]
command=kylin.sh start
priority=10
autostart=true
user=kylin
redirect_stderr=true
stdout_logfile=/var/log/kylin/kylin.out