# manager/supervisor.d/cesi.conf
[program:cesi]
command=python3 /usr/lib/cesi/cesi/run.py --config-file /etc/cesi.conf.toml
priority=10
autostart=true
user=root
redirect_stderr=true
stdout_logfile=/var/log/cesi/cesi.out