# spark/supervisor.d/spark.conf
[program:spark_livy]
command=livy-server
priority=10
autostart=true
user=spark
redirect_stderr=true
stdout_logfile=/var/log/livy/livy.out