# hbase/supervisor.d/hbase.conf
[program:hbase]
command=start-hbase.sh
priority=10
autostart=true
user=hbase
redirect_stderr=true
stdout_logfile=/var/log/hbase/hbase.out